# .crontab
# Minuto	Hora	Dia	Mes	Dia-Semana	Orden
0		*	*	*	*		backup-linuxcloud
